# server.cfg
[supervisord]
nodaemon=true
user=user-app

logfile=/dev/null
logfile_maxbytes=0

pidfile=/tmp/supervisord.pid

[program:nginx]
command=nginx -g 'daemon off;' -e /dev/stderr

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

stopsignal=QUIT
stopasgroup=true

stopwaitsecs=20
killasgroup=true

[program:app]
command=sh /etc/gunicorn/launcher.sh wsgi:app --config /etc/gunicorn/config.py

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

stopasgroup=false

stopwaitsecs=20
killasgroup=true
